# deploy/launchagents/com.ogz.brain-api.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.ogz.brain-api</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/python3</string>
    <string>/Users/abarihm/Desktop/ogz-knowledge/scripts/brain_api_launcher.py</string>
  </array>
  <key>KeepAlive</key><true/>
  <key>RunAtLoad</key><true/>
  <key>ThrottleInterval</key><integer>10</integer>
  <key>StandardOutPath</key><string>/Users/abarihm/logs/brain_api.out.log</string>
  <key>StandardErrorPath</key><string>/Users/abarihm/logs/brain_api.err.log</string>
</dict>
</plist>
